# dist/opt/retroroller/services/oga_events.service
[Unit]
Description=odroidgo2 evdev listener

[Service]
ExecStart=/usr/local/bin/oga_events.py
Type=simple

[Install]
WantedBy=multi-user.target
